# unit.service
[Unit]
After=network.target
Description=my service
After=libvirtd.service

[Service]
ExecStart=/path/to/start -f -cd --long-bool-flag
[Install]
